# clawd/dashboard/dashboard.service
[Unit]
Description=Liam's Project Dashboard
After=network.target

[Service]
Type=simple
User=liam
WorkingDirectory=/home/liam/clawd/dashboard
ExecStart=/usr/bin/python3 app.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
