# infra/bootstrap/openclawd.service
[Unit]
Description=OpenClaw Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=openclaw
Group=openclaw
TimeoutStartSec=20min

# State and config directories
Environment=OPENCLAW_STATE_DIR=/var/lib/openclaw
Environment=OPENCLAW_CONFIG_PATH=/var/lib/openclaw/openclaw.json

# Repository locations
# DEV_REPO: Where the agent reads source code (workspace, read-only for runtime)
# RUNTIME_DIR: Where the protected runtime lives (built artifacts)
Environment=OPENCLAW_REPO_DIR=/var/lib/openclaw/workspace/openclaw
Environment=OPENCLAW_RUNTIME_DIR=/opt/openclaw/runtime
Environment=OPENCLAW_REPO_URL=https://github.com/aron98/openclaw.git
Environment=OPENCLAW_REPO_REF=main

# Node configuration for control panel access
Environment=HOME=/var/lib/openclaw
Environment=NPM_CONFIG_PREFIX=/var/lib/openclaw/.npm-global
Environment=NODE_OPTIONS=--max-old-space-size=3072

# Load additional environment from file if present
EnvironmentFile=-/var/lib/openclaw/openclawd.env

# Update PATH to include user-local npm global bin
Environment=PATH=/var/lib/openclaw/.npm-global/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin

# Sync secrets from AWS Secrets Manager
ExecStartPre=/usr/local/bin/openclaw-sync-secrets arn:aws:secretsmanager:eu-central-1:548135083003:secret:openclaw/prod-nTa8N2

# Start from the PROTECTED runtime directory, not the dev repo
ExecStart=/usr/bin/env node /opt/openclaw/runtime/openclaw.mjs gateway --port 18789 --verbose

WorkingDirectory=/var/lib/openclaw
Restart=on-failure
RestartSec=5
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
